# root/etc/services.d/zoneminder/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# ZoneMinder
# Take down the S6 supervision tree when ZoneMinder fails to start/crashes
# ==============================================================================

s6-svscanctl -t /var/run/s6/services