# root/etc/services.d/zoneminder/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# ZoneMinder
# Runs ZoneMinder
# ==============================================================================

echo "Starting ZoneMinder..."
exec /usr/bin/zmpkg.pl start
