# root/etc/services.d/apache2/finish
#!/usr/bin/with-contenv bash

exec /usr/sbin/apachectl -k graceful-stop
