# root/etc/services.d/apache2/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Apache2
# Tell Apache2 to stop when container stops
# ==============================================================================

exec /usr/sbin/apache2 -k stop
