# root/etc/services.d/apache2/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Apache2
# Runs Apache2
# ==============================================================================

echo "[Apache2] Starting Apache2..."

source /etc/apache2/envvars
/usr/sbin/apache2ctl -D FOREGROUND | sed -u "s/^/[Apache2] /"
