# root/etc/services.d/mariadb-configure/run
#!/usr/bin/with-contenv bash
. "/usr/local/bin/logger"
# ==============================================================================
# mariadb-configure
# Configures db using /etc/mariadbconfigure.d
# ==============================================================================

# Reconfigure to be oneshot
/bin/s6-svc -O /var/run/s6/services/mariadb-configure

echo "Waiting for MariaDB to start" | info
# Wait for db to be up before configuring
/bin/s6-svwait -U /var/run/s6/services/mariadb

for f in /etc/mariadbconfigure.d/*.sh; do
  bash "$f"
done

echo "Upgrading db if necessary" | info
s6-setuidgid www-data /usr/bin/zmupdate.pl -nointeractive | info

echo "Refreshing db" | info
s6-setuidgid www-data /usr/bin/zmupdate.pl -nointeractive -f | info

# Notify s6 service is up
fdmove 1 3 printf "done\n"
