# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
. "/usr/local/bin/logger"
# ==============================================================================
# Nginx
# Runs Nginx
# ==============================================================================

echo "Starting Nginx..."

exec nginx -g 'daemon off;'
